# roles/mount/templates/music.automount
{{ ansible_managed }}

[Unit]
Description=Automount music storage

[Automount]
Where={{ mount_directory }}

[Install]
WantedBy=multi-user.target
